# u1/systemd/video-rx.service
# /etc/systemd/system/video-rx.service (на У1)

[Unit]
Description=Video RX (H.264 RTP -> HDMI via cage+waylandsink)
After=network-online.target
Wants=network-online.target
Conflicts=display-manager.service

[Service]
Type=simple
ExecStart=/opt/u1u2-bridge/u1/video_rx.sh
Restart=always
RestartSec=3
Nice=-5
Environment=GST_DEBUG=2

[Install]
WantedBy=multi-user.target
