# u2/systemd/video-tx.service
# /etc/systemd/system/video-tx.service (на У2)

[Unit]
Description=Video TX (VRX -> H.264 RTP -> CPE710)
After=network-online.target dev-video0.device
Wants=network-online.target

[Service]
Type=simple
# Видео-peer (PEER_HOST) — из env, симметрично CRSF. '-' = отсутствие файла
# не валит юнит (в tunnel файла нет, video_tx.sh берёт fallback 10.8.0.6).
EnvironmentFile=-/etc/u1u2-bridge/video.env
ExecStart=/opt/u1u2-bridge/u2/video_tx.sh
Restart=always
RestartSec=3
Nice=-5

[Install]
WantedBy=multi-user.target
